# services/lakshmi/org.pancia.lakshmi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>org.pancia.lakshmi</string>
        <key>ProgramArguments</key>
        <array>
            <string>/Users/anthony/dotfiles/services/lakshmi/dotfiles-services-lakshmi.sh</string>
        </array>
        <key>RunAtLoad</key>
        <true/>
        <key>KeepAlive</key>
        <true/>
        <key>StandardOutPath</key>
        <string>/Users/anthony/.log/services/lakshmi.log</string>
        <key>StandardErrorPath</key>
        <string>/Users/anthony/.log/services/lakshmi.log</string>
    </dict>
</plist>
